# supervisord/supervisor/log.conf
[program:log]
command=sh /script/log.sh

